# luci/themes/luci-theme-ssn/root/etc/uci-defaults/30_luci-theme-ssn
#!/bin/sh
uci batch <<-EOF
	set luci.themes.SSN=/luci-static/ssn
	set luci.main.mediaurlbase=/luci-static/ssn
	commit luci
EOF
exit 0
